# etc/profile.d/ant.sh
export ANT_OPTS="-Dhttp.proxyHost=proxy.kuins.net -Dhttp.proxyPort=8080"
