# etc/profile.d/proxy.sh
export http_proxy=http://proxy.kuins.net:8080/
export ftp_proxy=http://proxy.kuins.net:8080/
export HTTP_PROXY=http://proxy.kuins.net:8080/
export FTP_PROXY=http://proxy.kuins.net:8080/
